# mqtt_sub.service
[Unit]
Description=MQTT Subscribe
After=multi-user.target
[Service]
Type=simple
Restart=always
ExecStart=/usr/bin/python3 /home/ubuntu/app/shoestring-hackathon-2022/mqtt_subscribe.py
[Install]
WantedBy=multi-user.target